# skills/google-workspace/scripts/workspace-token-maintainer.timer
[Unit]
Description=Run Google Workspace Token Maintenance every 45 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=45min
Persistent=true

[Install]
WantedBy=timers.target
