# skills/google-workspace/scripts/workspace-sync.service
[Unit]
Description=Google Workspace Cache Sync
After=network-online.target
Wants=network-online.target
# Only run on AC power
ConditionACPower=true

[Service]
Type=oneshot
# Point cache to the big NVMe drive
Environment=GOOGLE_WORKSPACE_DIR=%h/.google_workspace
Environment=WORKSPACE_CACHE_DIR=/mnt/ai-data/google-workspace-cache
ExecStartPre=/usr/bin/mkdir -p /mnt/ai-data/google-workspace-cache
ExecStart=/usr/bin/uv run %h/.agents/skills/google-workspace/scripts/cache.py sync
TimeoutStartSec=300
# Don't mark as failed on network errors — just skip this cycle
Restart=no
SuccessExitStatus=0 1
StandardOutput=journal
StandardError=journal
