# skills/google-workspace/scripts/workspace-token-maintainer.service.template
[Unit]
Description=Google Workspace Token Maintenance
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%UV_PATH% run %INSTALL_DIR%/scripts/maintain_token.py
StandardOutput=journal
StandardError=journal
